# ezcommon-frontend/systemd/ezcommon-frontend.service
[Unit]
Description=EZCommon Frontend (Next.js production)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ec2-user/full-web-app-1123/ezcommon-frontend
Environment="PORT=3000"
Environment="HOSTNAME=0.0.0.0"
ExecStart=/home/ec2-user/.nvm/versions/node/v20.19.6/bin/node /home/ec2-user/full-web-app-1123/ezcommon-frontend/node_modules/.bin/next start --hostname 0.0.0.0 --port ${PORT}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
